# rc.local
#!/bin/bash

cd /opt/nvidia/deepstream/deepstream/sources/apps/Census-Deepstream
./start.sh& 

